# deploy/crontab
0 4 * * 3 /vol/vabusinesses.org/htdocs/scripts/update.sh
